# supervisor.conf
[program:xr]
command=/home/src/extinctionr/run.sh
directory=/home/src/extinctionr/
environment=PATH="/bin:/usr/bin:/usr/local/bin:/home/dcs/.local/bin/"
user=dcs
pidfile=/var/run/xr.pid
