# supervisor.conf
[program:xr]
command=/home/src/extinctionr/run.sh
directory=/home/src/extinctionr/
